# reminder-bot.service
[Unit]
Description=Telegram Reminder Bot
After=network.target
Wants=network.target

[Service]
Type=simple
User=botuser
Group=botuser
WorkingDirectory=/opt/reminder-bot
Environment=PYTHONPATH=/opt/reminder-bot
ExecStart=/opt/reminder-bot/venv/bin/python /opt/reminder-bot/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Переменные окружения
EnvironmentFile=/opt/reminder-bot/.env

# Ограничения ресурсов
MemoryLimit=256M
CPUQuota=50%

# Безопасность
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/reminder-bot

[Install]
WantedBy=multi-user.target
